# cmdb_full_inventory.json
{
  "queries": {
    "cmdb_software": {
      "query": "SELECT name, version, publisher, install_date FROM programs;",
      "interval": 3600,
      "description": "Installed software inventory"
    },
    "cmdb_network": {
      "query": "SELECT interface, address, mask FROM interface_addresses WHERE address NOT IN ('127.0.0.1', '::1', '0.0.0.0');",
      "interval": 600,
      "description": "Network interfaces inventory"
    },
    "cmdb_usb_devices": {
      "query": "SELECT vendor, model, serial, usb_port FROM usb_devices;",
      "interval": 600,
      "description": "Connected USB devices"
    },
    "cmdb_users": {
      "query": "SELECT username, description, directory, uuid, type FROM users;",
      "interval": 3600,
      "description": "Local users list"
    },
    "cmdb_patches": {
      "query": "SELECT hotfix_id, description, installed_by, installed_on FROM patches;",
      "interval": 7200,
      "description": "Windows Updates / Patches"
    },
    "cmdb_startup_items": {
      "query": "SELECT name, path, type, status FROM startup_items;",
      "interval": 3600,
      "description": "Autorun items"
    },
    "cmdb_logical_drives": {
      "query": "SELECT device_id, free_space, size, file_system FROM logical_drives;",
      "interval": 3600,
      "description": "Disk space information"
    }
  }
}
