# pack_software_audit.json
{
  "queries": {
    "software_installed_programs": {
      "query": "SELECT name, version, publisher, install_date FROM programs;",
      "interval": 14400,
      "platform": "windows",
      "description": "Сбор списка всех установленных программ Windows (аналог Программы и компоненты)"
    },
    "software_startup_items": {
      "query": "SELECT name, path, args, source FROM startup_items;",
      "interval": 14400,
      "platform": "windows",
      "description": "Отслеживание программ, прописанных в автозагрузку (реестр, папки Startup)"
    }
  }
}
